# rootfs/etc/services.d/apache2/run
#!/bin/sh

exec /usr/local/bin/docker-php-entrypoint /usr/local/bin/apache2-foreground